# packaging/botka.service
[Unit]
Description=Botka - AI Task Scheduler & Chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
EnvironmentFile=/home/pi/projects/botka/.env
Environment=PATH=/home/pi/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/bin/botka
WorkingDirectory=/home/pi/projects/botka
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=botka

[Install]
WantedBy=multi-user.target
